# redmine_backup.cron
#!/usr/bin/env bash

BACKUP="$HOME/redmine_backup/redmine_backup.sh"
DOM=$(date +%d);
DOW=$(date +%a);
TYPE="d";

if   [ "${DOM}" == "01"  ]; then TYPE="m";
elif [ "${DOW}" == "Sun" ]; then TYPE="w";
fi;

$BACKUP ${TYPE} > /dev/null 2>&1

exit 0;